# systemd/hackrf-scanner.service
[Unit]
Description=HackRF Scanner - Real-time Spectrum Analyzer
After=network.target
Wants=openwebrx.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/projects/stinkster/hackrf
ExecStart=/home/pi/projects/stinkster/hackrf/venv/bin/python /home/pi/projects/stinkster/hackrf/spectrum_analyzer.py
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=8090
Environment=OPENWEBRX_URL=ws://localhost:8073/ws/
Environment=SAMPLE_RATE=2400000
Environment=CENTER_FREQ=145000000
Environment=FFT_SIZE=1024

# Resource limits for Raspberry Pi
MemoryLimit=256M
CPUQuota=25%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hackrf-scanner

[Install]
WantedBy=multi-user.target
