# systemd/kismet-operations-center.service
[Unit]
Description=Kismet Operations Center - Node.js WiFi Monitoring Dashboard
Documentation=https://github.com/yourusername/stinkster
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/projects/stinkster_malone/stinkster/src/nodejs/kismet-operations
Environment="NODE_ENV=production"
Environment="PATH=/usr/bin:/usr/local/bin:/home/pi/.nvm/versions/node/v18.17.0/bin"
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/kismet-operations-center.log
StandardError=append:/var/log/kismet-operations-center.log

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/projects/stinkster_malone/stinkster/src/nodejs/kismet-operations
ReadWritePaths=/var/log
ReadWritePaths=/home/pi/tmp
ReadWritePaths=/home/pi/kismet_ops
ReadWritePaths=/home/pi/WigletoTAK

[Install]
WantedBy=multi-user.target