# outputs/webhook_port_20250618_233716/final/webhook.service
[Unit]
Description=Stinkster Webhook Service
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/stinkster/src/nodejs/webhook-service
Environment="NODE_ENV=production"
Environment="WEBHOOK_PORT=8002"
ExecStart=/usr/bin/node webhook.js
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webhook-service

[Install]
WantedBy=multi-user.target